# deploy/oem/root/etc/systemd/system/suse-assist-firstboot.service
# One-shot firstboot job: populate the assistant data volume before the main
# service starts. If /usr/share/suse-assist/offline-bundle.tar.gz exists, it is
# imported. Otherwise this falls back to networked image pull + doc ingestion.

[Unit]
Description=openSUSE AI Onboarding — first boot model & index setup
After=network-online.target podman.socket
Wants=network-online.target
Before=suse-assist.service
ConditionPathExists=!/var/lib/suse-assist/.firstboot-done

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=3600
Environment=SUSE_ASSIST_IMAGE=ghcr.io/anujagrawal380/opensuse-leap-ai-startup-guide:latest
Environment=SUSE_ASSIST_BUNDLE=/usr/share/suse-assist/offline-bundle.tar.gz
ExecStart=/usr/libexec/suse-assist-firstboot

[Install]
WantedBy=multi-user.target
